# hub/systemd/church-display-media-index.timer
[Unit]
Description=Refresh Church Display media metadata index periodically

[Timer]
OnBootSec=45s
OnUnitActiveSec=2min
Persistent=true
RandomizedDelaySec=15s

[Install]
WantedBy=timers.target
